# cookbooks/mongodb/files/default/mongodb
#!/sbin/runscript
# Copyright 1999-2009 Gentoo Foundation
# Distributed under the terms of the GNU General Public License v2
# $Header: $

depend() {
  need net
}

MONGODB_EXEC="${MONGODB_EXEC:-/usr/local/mongodb/bin/mongod}"
MONGODB_MASTER_DATA="${MONGODB_MASTER_DATA:-/db/mongodb/masterdb}"
MONGODB_SLAVE_DATA="${MONGODB_SLAVE_DATA:-/db/mongodb/slavedb}"
MONGODB_USER="${MONGODB_USER:-mongodb}"
MONGODB_MASTER_PID_FILE="${MONGODB_MASTER_PID_FILE:-/var/run/mongodb/master.pid}"
MONGODB_SLAVE_PID_FILE="${MONGODB_SLAVE_PID_FILE:-/var/run/mongodb/slave.pid}"

MONGODB_IP="${MONGODB_IP:-127.0.0.1}"
MONGODB_MASTER_PORT="${MONGODB_MASTER_PORT:-27017}"
MONGODB_SLAVE_PORT="${MONGODB_SLAVE_PORT:-27018}"

start() {
  ebegin "Starting Master ${SVCNAME}"

  args="--dbpath $MONGODB_MASTER_DATA --port $MONGODB_MASTER_PORT --bind_ip $MONGODB_IP"

  start-stop-daemon --background --start --chuid ${MONGODB_USER} --pidfile \
  "${MONGODB_MASTER_PID_FILE}" --make-pidfile --exec ${MONGODB_EXEC} \
  -- ${args} run
  eend $?

  ebegin "Starting Slave ${SVCNAME}"

  args="--dbpath $MONGODB_SLAVE_DATA --port $MONGODB_SLAVE_PORT --bind_ip $MONGODB_IP"

  start-stop-daemon --background --start --chuid ${MONGODB_USER} --pidfile \
  "${MONGODB_SLAVE_PID_FILE}" --make-pidfile --exec ${MONGODB_EXEC} \
  -- ${args} run
  eend $?
}

stop() {
  ebegin "Stopping Master ${SVCNAME}"
  start-stop-daemon --stop --exec ${MONGODB_EXEC} --pidfile ${MONGODB_MASTER_PID_FILE}
  eend $? "Failed to stop Master ${SVCNAME}"

  ebegin "Stopping Slave ${SVCNAME}"
  start-stop-daemon --stop --exec ${MONGODB_EXEC} --pidfile ${MONGODB_SLAVE_PID_FILE}
  eend $? "Failed to stop Slave ${SVCNAME}"
}

# vim:set filetype=sh:
